# packaging/locallexis-hub.service
# systemd unit for the headless LocalLexis hub.
#
# Install:
#   sudo cp packaging/locallexis-hub.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now locallexis-hub
#
# Adjust User, WorkingDirectory, and the venv path to match the host
# you're running on. The service expects the speechtotext package
# installed into the referenced venv with the [api] extra:
#   /opt/locallexis/.venv/bin/pip install -e /opt/locallexis/SpeechToText[api]

[Unit]
Description=LocalLexis hub (headless)
Documentation=https://github.com/lieuwe89/LocalLexis
After=network.target

[Service]
Type=exec
User=locallexis
Group=locallexis
WorkingDirectory=/opt/locallexis
Environment="LOCALLEXIS_HOST=0.0.0.0"
Environment="LOCALLEXIS_PORT=8765"
# Self-signed cert + key are generated lazily under the service user's
# data dir on first start. Override to 0 only if TLS is terminated by a
# reverse proxy in front of this unit.
Environment="LOCALLEXIS_TLS_ENABLED=1"
# REQUIRED on this public (0.0.0.0) bind: without an admin token the hub
# refuses to start, because the admin API (config, transcripts, pairing-token
# minting, job control) would otherwise be unauthenticated on the network.
# Generate a secret — `openssl rand -hex 32` — and set it below (or point an
# EnvironmentFile= at one). To bind without a token on a trusted/isolated
# network, set LOCALLEXIS_ALLOW_ANONYMOUS=1 instead.
#Environment="LOCALLEXIS_API_TOKEN=replace-with-a-generated-secret"
ExecStart=/opt/locallexis/.venv/bin/locallexis-hub
Restart=on-failure
RestartSec=5
# Confine the process: it only needs read of code + write of its data
# dir (the user's HOME by default; LOCALLEXIS_DATA_DIR will land in a
# future commit).
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=read-only
ReadWritePaths=/var/lib/locallexis

[Install]
WantedBy=multi-user.target
